# systemv.socket
[Socket]
ListenStream=/var/run/systemv.sk

[Install]
WantedBy=sockets.target
