# services/systemv.service
[Service]
ExecStart=/usr/bin/systemv

[Install]
WantedBy=multi-user.agent
